# deploy/breachnotifier.timer
[Unit]
Description=Run BreachNotifier daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
